# crontab
* * * * * /bin/bash -l -c 'cd /usr/src/app && RAILS_ENV=development bundle exec rake punch:matching --silent'